# deploy/systemd/ton-agent-bot.service.template
[Unit]
Description=TON Agent Ads Bot
After=network.target ton-agent-api.service
Wants=ton-agent-api.service

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_ROOT__/current/apps/bot
Environment=DOTENV_CONFIG_PATH=__APP_ROOT__/shared/.env
ExecStart=/usr/bin/env node dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
